# daemon_assets/reme.service
[Unit]
Description=Reme

[Service]
ExecStart=/usr/bin/reme -daemon

[Install]
WantedBy=multi-user.target
